# etc/host-tools.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/sbin/host-tools notify-boot

[Install]
WantedBy=default.target
