# etc/systemd/host-tools-alive.service
[Unit]
Description=Keepalive notification

[Service]
ExecStart=-/usr/local/sbin/host-tools notify-alive

[Install]
WantedBy=host-tools-alive.timer
